# cron/weatherStationTimeTable.crontab
SHELL=/bin/sh

# min  hour  day  month  weekday  command
*/30   *     *    *      *        sh /home/pi/station/weatherstation-maincontrolcenter/programs/callWeatherStation.sh > /dev/null 2>&1
15     5     *    *      *        sh /home/pi/station/weatherstation-maincontrolcenter/bin/backupData.sh > /dev/null 2>&1
